# deploy/ws_sink/ws_sink-stop.service
[Unit]
Description=Stop ws_sink (helper for ws_sink-stop.timer)

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop ws_sink.service
